# systemd/mjpg-streamer.service
[Unit]
Description=mjpg-streamer (MJPEG webcam server)
After=network.target

[Service]
User=maker2
Group=video
EnvironmentFile=/etc/default/mjpg-streamer
ExecStartPre=/bin/sh -c 'test -x /usr/local/bin/mjpg_streamer || /bin/false'
ExecStart=/bin/sh -c 'eval "/usr/local/bin/mjpg_streamer $MJPG_INPUT_OPTS $MJPG_OUTPUT_OPTS"'
Restart=on-failure
RestartSec=3
# Ensure service doesn't run until video device permissions are correct
# Adjust user/group (default 'pi' and 'video') as needed for your system.

[Install]
WantedBy=multi-user.target